# linux-c7-libcanberra/pkg-plist
usr/bin/canberra-boot
usr/lib/libcanberra-0.30/libcanberra-alsa.so
usr/lib/libcanberra-0.30/libcanberra-gstreamer.so
usr/lib/libcanberra-0.30/libcanberra-multi.so
usr/lib/libcanberra-0.30/libcanberra-null.so
usr/lib/libcanberra-0.30/libcanberra-pulse.so
usr/lib/libcanberra.so.0
usr/lib/libcanberra.so.0.2.5
usr/lib/systemd/system/canberra-system-bootup.service
usr/lib/systemd/system/canberra-system-shutdown-reboot.service
usr/lib/systemd/system/canberra-system-shutdown.service
usr/lib64/libcanberra-0.30/libcanberra-alsa.so
usr/lib64/libcanberra-0.30/libcanberra-gstreamer.so
usr/lib64/libcanberra-0.30/libcanberra-multi.so
usr/lib64/libcanberra-0.30/libcanberra-null.so
usr/lib64/libcanberra-0.30/libcanberra-pulse.so
usr/lib64/libcanberra.so.0
usr/lib64/libcanberra.so.0.2.5
%%PORTDOCS%%%%DOCSDIR%%/LGPL
%%PORTDOCS%%%%DOCSDIR%%/README

